# rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi
su -l pi -c "export LD_LIBRARY_PATH=~/xmpp-remote-agent/out/lib;~/xmpp-remote-agent/out/bin/bboxsmsrv --boardtype=RASPI_BPLUS" &
su -l pi -c "export LD_LIBRARY_PATH=~/xmpp-remote-agent/out/lib;~/xmpp-remote-agent/out/bin/sysmgr --boardtype=RASPI_BPLUS" &
su -l pi -c "export LD_LIBRARY_PATH=~/xmpp-remote-agent/out/lib;~/xmpp-remote-agent/out/bin/dispsrv --iface=eth0" &
su -l pi -c "export LD_LIBRARY_PATH=~/xmpp-remote-agent/out/lib;~/xmpp-remote-agent/out/bin/xmproxysrv --iface=eth0 --boardtype=RASPI_BPLUS --loginfile=/home/pi/xmpp-remote-agent/out/etc/xmproxy/xmpp-login.txt --aliaslist=/home/pi/xmpp-remote-agent/out/etc/xmproxy/xmpp-alias-list.txt --botname=/home/pi/xmpp-remote-agent/out/etc/xmproxy/xmpp-bot-name.txt --usbgsm=1" &
exit 0
